# rabbitmq/supervisord.conf
[supervisord]
nodaemon=true

[program:cluster]
command=/make_cluster.sh

[program:rabbit_server]
command=/bin/bash -c "/docker-entrypoint.sh rabbitmq-server"
